# coint4/scripts/optimization/systemd/autonomous-confirm-diversity-guard-agent.service
[Unit]
Description=Autonomous Confirm Diversity Guard Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=__ROOT_DIR__
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/env python3 __ROOT_DIR__/scripts/optimization/confirm_diversity_guard_agent.py

[Install]
WantedBy=default.target
