# coint4/scripts/optimization/systemd/autonomous-confirm-dispatch-agent.timer
[Unit]
Description=Run Autonomous WFA Confirm Dispatch Agent every 1 minute

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=15s
Persistent=true
Unit=autonomous-confirm-dispatch-agent.service

[Install]
WantedBy=timers.target
