# coint4/scripts/optimization/systemd/autonomous-progress-guard-agent.timer
[Unit]
Description=Run Autonomous Progress Guard every 3 minutes

[Timer]
OnBootSec=90s
OnUnitActiveSec=3min
AccuracySec=20s
Persistent=true
Unit=autonomous-progress-guard-agent.service

[Install]
WantedBy=timers.target
