# coint4/scripts/optimization/systemd/autonomous-single-writer-guard-agent.timer
[Unit]
Description=Run Single Writer Guard every 5 minutes

[Timer]
OnBootSec=140s
OnUnitActiveSec=5min
AccuracySec=20s
Persistent=true
Unit=autonomous-single-writer-guard-agent.service

[Install]
WantedBy=timers.target
