# coint4/scripts/optimization/systemd/autonomous-confirm-diversity-guard-agent.timer
[Unit]
Description=Run Confirm Diversity Guard every 2 minutes

[Timer]
OnBootSec=110s
OnUnitActiveSec=2min
AccuracySec=15s
Persistent=true
Unit=autonomous-confirm-diversity-guard-agent.service

[Install]
WantedBy=timers.target
