# coint4/scripts/optimization/systemd/autonomous-deterministic-error-blacklist-agent.timer
[Unit]
Description=Run Deterministic Error Blacklist Agent every 3 minutes

[Timer]
OnBootSec=120s
OnUnitActiveSec=3min
AccuracySec=20s
Persistent=true
Unit=autonomous-deterministic-error-blacklist-agent.service

[Install]
WantedBy=timers.target
