# coint4/scripts/optimization/systemd/autonomous-wfa-watchdog.timer
[Unit]
Description=Run Autonomous WFA Watchdog every 3 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=3min
AccuracySec=15s
Persistent=true
Unit=autonomous-wfa-watchdog.service

[Install]
WantedBy=timers.target
