# coint4/scripts/optimization/systemd/autonomous-promotion-gatekeeper-agent.timer
[Unit]
Description=Run Autonomous Promotion Gatekeeper Agent every 1 minute

[Timer]
OnBootSec=90s
OnUnitActiveSec=1min
AccuracySec=10s
Persistent=true
Unit=autonomous-promotion-gatekeeper-agent.service

[Install]
WantedBy=timers.target

