# coint4/scripts/optimization/systemd/autonomous-search-director-agent.service
[Unit]
Description=Autonomous Search Director Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=__ROOT_DIR__
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/env python3 __ROOT_DIR__/scripts/optimization/search_director_agent.py

[Install]
WantedBy=default.target

