# coint4/scripts/optimization/systemd/autonomous-backlog-janitor-agent.timer
[Unit]
Description=Run Autonomous Backlog Janitor every 2 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=2min
RandomizedDelaySec=20s
Unit=autonomous-backlog-janitor-agent.service

[Install]
WantedBy=timers.target
