# coint4/scripts/optimization/systemd/autonomous-queue-seeder.timer
[Unit]
Description=Run Autonomous WFA Queue Seeder every 2 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=2min
AccuracySec=15s
Persistent=true
Unit=autonomous-queue-seeder.service

[Install]
WantedBy=timers.target
