# coint4/scripts/optimization/systemd/ops-sentinel-agent.timer
[Unit]
Description=Run Ops Sentinel Agent every 2 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=2min
AccuracySec=15s
Persistent=true
Unit=ops-sentinel-agent.service

[Install]
WantedBy=timers.target
